# superkeg.service
[Unit]
Description=Superkeg Beer Management System (Flask + Flow Meters)
After=network.target
Wants=network.target

[Service]
Type=simple
User=pi
Group=pi
WorkingDirectory=/home/pi/superkeg
ExecStart=/usr/bin/python /home/pi/superkeg/start_superkeg.py
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal

# Environment variables
Environment=PYTHONPATH=/home/pi/superkeg
Environment=FLASK_ENV=production

# Security settings
NoNewPrivileges=true
PrivateTmp=true

[Install]
WantedBy=multi-user.target